# pcbaby/supervisor.ini
[program:pcbaby_ask]
command=scrapy crawl pcbabay_ask
directory=/root/kst/pcbaby
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=5
stdout_logfile=/root/kst/tuniu/crawlstd.log
redirect_stderr=true
