# askspider/supervisor.ini
[program:huicaiba_ask]
command=scrapy crawl huicaiba_ask 
directory=/home/heyong/kst/askspider
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=5
stdout_logfile=/home/heyong/kst/askspider/log/huicaiba.log
redirect_stderr=true
