# tuniu/supervisor.ini
[program:tuniu_ask]
command=scrapy crawl tuniu_ask 
directory=/home/heyong/kst/tuniu
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=5
stdout_logfile=/home/heyong/kst/tuniu/crawl.log
redirect_stderr=true
