# tuniu/supervisor.ini
[program:tuniu_ask]
command=scrapy crawl tuniu_ask
directory=/root/kst/tuniu
autostart=false
autorestart=false
startsecs=3
stopwaitsecs=5
stdout_logfile=/root/kst/tuniu/crawlstd.log
redirect_stderr=true
